# frontend/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:next]
command=/usr/local/bin/yarn start:prod
directory=/app
